# rootfs/etc/supervisor/conf.d/spigot.conf
[program:spigot]
command=/my_service minecraft_server
redirect_stderr=true
autostart=true
stdout_syslog=true

