# rootfs/etc/supervisor/conf.d/spigot.conf
[program:spigot]
command=/waitjava.sh minecraft_server spigot.jar
redirect_stderr=true
autostart=true
stdout_syslog=true

